# rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Starts app
# ==============================================================================
exec system-bridge
