# rootfs/etc/services.d/app/run
#!/bin/bash
# ==============================================================================
# Starts app
# ==============================================================================
cd /usr/local/share/system-bridge || exit 1

export SB_APPDATA=/data/system-bridge
export SB_CLI=true
export SB_CONTAINER=true
export SB_TRAY=false

exec system-bridge
